# resources/nodejs-agent.service
[Unit]


[Service]
ExecStart=/usr/bin/node /home/eitan/projects/nodejs-agent/index.js

[Install]
